# imports/archillx-evidence/deploy/systemd/archillx-audit-archive.service
[Unit]
Description=ArcHillx audit archive job
After=network-online.target archillx.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/archillx
EnvironmentFile=/opt/archillx/.env.prod
ExecStart=/bin/sh -lc './scripts/archive_audit.sh'
