# imports/archillx-evidence/deploy/systemd/archillx-restore-drill.timer
[Unit]
Description=ArcHillx Monthly Restore Drill Timer

[Timer]
OnCalendar=monthly
Persistent=true
Unit=archillx-restore-drill.service

[Install]
WantedBy=timers.target
